# com.futureface.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.futureface.autopush</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ahmedebeid/Documents/futureface-app/autopush.sh</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/Users/ahmedebeid/Documents/futureface-app/.push_trigger</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/futureface-push.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/futureface-push-error.log</string>
</dict>
</plist>
